# deployment/dewared.service.tmpl
[Unit]
Description=deware daemon
After=multi-user.target

[Service]
ExecStart=${python} ${install_location}bin/dewared
Type=simple
Restart=always

[Install]
WantedBy=default.target
